# filesystem/inittab
::sysinit:/etc/init.d/rcS
::respawn:/bin/runsvdir /etc/rc.d
tty1::respawn:/sbin/getty 38400 tty1
tty2::respawn:/sbin/getty 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3
::respawn:/sbin/getty -L ttyS0 9600 vt100
::ctrlaltdel:/bin/umount -a -r
::shutdown:/usr/bin/umount -a -r
