# filesystem/inittab
::sysinit:/sbin/systemctl init
::askfirst:-/bin/sh
::ctrlaltdel:/bin/umount -a -r
::shutdown:/usr/bin/umount -a -r
