# filesystem/inittab
::sysinit:/etc/init.d/rcS
::respawn:/bin/runsvdir /etc/rc.d
::askfirst:-/bin/sh
::ctrlaltdel:/bin/umount -a -r
::shutdown:/usr/bin/umount -a -r
